# periodic/monthly/backup-monthly.sh
#! /usr/sh
python /app/main.py -i monthly